# server/keepalive.service
# Counters OCI Always Free idle reclaim (reclaimed if 7-day 95th-percentile
# CPU/network/memory < 20%). Generates ~30% total CPU (0.6 of 2 OCPUs) for 2h daily:
# 14h/week > required 8.4h (5% of 168h), keeping the 95th-percentile CPU >= 20%.
# Rationale: docs/superpowers/specs/2026-07-11-oci-claude-server-design.md
[Unit]
Description=CPU keep-alive against OCI idle reclaim

[Service]
Type=oneshot
ExecStart=/usr/bin/stress-ng --cpu 1 --cpu-load 60 --timeout 7200
Nice=19
